# supervisor_app.conf
[include]
files=/etc/supervisor/supervisord.conf

[program:awesome_app]
#directory=/mnt/c/my_study/docker/flask_web
directory=./
command=gunicorn -w 3 -b 0.0.0.0:8000 myapp:app